# wsl-final-xfce/scripts/display-monitor.service
[Unit]
Description=Display Mode Monitor Watchdog
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash %h/scripts/display-monitor.sh
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=display-monitor

[Install]
WantedBy=default.target
